# gearman-worker/supervisord.conf
[supervisord]
nodaemon=true

[program:mod_gearman_worker]
command=/bin/bash -c "/usr/bin/mod_gearman_worker --config=/etc/mod_gearman/worker.conf"
autostart=true
autorestart=true
redirect_stderr=true